# config/foreman/systemd/process.service.erb
# Based on foreman-0.88.1/data/export/systemd/process.service.erb
[Unit]
PartOf=<%= app %>.target
StopWhenUnneeded=yes

[Service]
User=<%= user %>
WorkingDirectory=<%= engine.root %>
Environment=PORT=<%= port %>
Environment=PS=<%= process_name %>
<% engine.env.each_pair do |var,env| -%>
Environment="<%= var %>=<%= env %>"
<% end -%>
ExecStart=/bin/bash -lc 'exec -a "<%= app %>-<%= process_name %>" <%= process.command %>'
Restart=always
RestartSec=14s
StandardInput=null
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n
KillMode=mixed
TimeoutStopSec=<%= engine.options[:timeout] %>
# Limit memory usage to 500MB - it normally starts around 250M and slowly climbs from there
MemoryMax=500M
